# crontab
0 * * * * restart_caddy.sh
@reboot startsbots.sh
